# systemd/msmtp-queue.path
[Unit]
Description=Flush msmtp queue when a message is queued.

[Path]
PathExistsGlob=%D/msmtp/queue/*.mail

[Install]
WantedBy=paths.target
